# supervisord.conf
[supervisord]
nodaemon=true

[program:mumble]
command=/usr/sbin/murmurd -fg -ini /mumble-server.ini
autorestart=false